# embodied-claude/rootfs/etc/s6-overlay/s6-rc.d/init/run
#!/command/with-contenv bashio
# ==============================================================================
# Initialize Embodied Claude add-on
# ==============================================================================

bashio::log.info "Initializing Embodied Claude..."

# Validate configuration
if ! bashio::config.has_value 'tapo_camera_host'; then
    bashio::log.warning "Tapo camera host not configured yet"
fi

if ! bashio::config.has_value 'anthropic_api_key'; then
    bashio::log.info "No API key set - use 'claude login' in the web terminal to authenticate"
fi

bashio::log.info "Configuration validated"
